# deploy/protracklite-effort-reminder.timer
[Unit]
Description=Run ProtrackLite weekly effort reminder every Friday at 3 PM

[Timer]
OnCalendar=Fri 15:00
Timezone=Asia/Kolkata
Persistent=true
Unit=protracklite-effort-reminder.service

[Install]
WantedBy=timers.target
